# cloudflared.service
[Unit]
Description=Cloudflare Tunnel for JackStatz
After=network.target
Wants=network.target

[Service]
Type=simple
User=your-username
Group=your-username
WorkingDirectory=/path/to/jackstatz
ExecStart=/usr/local/bin/cloudflared tunnel --config /path/to/jackstatz/cloudflared.yml run jackstatz-tunnel
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared-jackstatz

# Security settings
NoNewPrivileges=yes
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/path/to/jackstatz
ReadWritePaths=/home/your-username/.cloudflared
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
